# my-app/src/python/etc/profile.d/conda.csh
version https://git-lfs.github.com/spec/v1
oid sha256:be53e3a062560faf43a360b8f3bd81a7e11d70b4f6219b8de902121c1b740d30
size 2436
