# my-app/src/python/etc/profile.d/conda.sh
version https://git-lfs.github.com/spec/v1
oid sha256:1e48c2efa96c33c17e0b9e67ddf0288566aca8243be73dba9ac459d39d561cb1
size 2219
