# my-app/src/python/pkgs/conda-23.3.1-py310haa95532_0/etc/profile.d/conda.csh
version https://git-lfs.github.com/spec/v1
oid sha256:365513ebdac197e668b275606d3fb82259acb54a8f9ee67ea518e5df288d7026
size 2460
