# my-app/src/python/pkgs/conda-23.3.1-py310haa95532_0/etc/profile.d/conda.sh
version https://git-lfs.github.com/spec/v1
oid sha256:bfeaa2cb7aaa896984a49bcaaf5859d9553c9712e5564f05dd39b98ee9ab00d6
size 2231
